# nodediag.service
[Unit]
Description=Nodediag
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/nodediag -s -v

[Install]
WantedBy=multi-user.target
